# pocketrugby-webhook.service
[Unit]
Description=PocketRugbyNL webhook deploy server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/PocketRugbyNL
ExecStart=node /home/pi/PocketRugbyNL/webhook-server.js
Restart=on-failure
RestartSec=5
EnvironmentFile=/etc/pocketrugby-webhook.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pocketrugby-webhook

[Install]
WantedBy=multi-user.target
